# freeserver-main.service
[Unit]
Description=My Free Server main Bot Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /root/telegrambot/free-server/main.py
Restart=always
RestartSec=10
User=root
Group=root

[Install]
WantedBy=multi-user.target
