# qr.service
[Unit]
Description=My QR Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /root/telegrambot/qr.py
Restart=always
RestartSec=10
User=root
Group=root

[Install]
WantedBy=multi-user.target
